# deployed-service-watcher.path
[Path]
PathModified=/srv/git/<++>/<++>/<++>.git/./index.lock

[Install]
WantedBy=multi-user.target
